# task6.service
[Unit]
Description=Logs
[Service]
User=root
ExecStart=/home/ubuntu/task6.sh
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=multi-user.target
